# rootfs/etc/s6/weewx/run
#!/bin/sh

exec doas -u weewx -- /usr/bin/weewxd --config /etc/weewx/weewx.conf
